# deployment/strava-backend.service
[Unit]
Description=Strava Activity Copilot Backend
After=network.target

[Service]
User=andy
WorkingDirectory=/home/andy/Dev/strava-activity-copilot
EnvironmentFile=/home/andy/Dev/strava-activity-copilot/backend/.env
ExecStart=/home/andy/Dev/strava-activity-copilot/venv/bin/uvicorn backend.main:app --host 127.0.0.1 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target
