# deployment/strava-mcp.service
[Unit]
Description=Strava MCP HTTP Server
After=network.target

[Service]
User=andy
WorkingDirectory=/var/www/strava-insight-portal
EnvironmentFile=/var/www/strava-insight-portal/backend/.env
# Note: MCP server needs some env vars too, specifically credentials if using refresh mode, 
# but our modified version relies on headers. However, it still uses client_id for refresh logic internally.
ExecStart=/var/www/strava-insight-portal/venv/bin/python mcp-server/src/strava_http_server.py
Restart=always

[Install]
WantedBy=multi-user.target
