# deployment/strava-backend.service
[Unit]
Description=Strava Insight Portal Backend
After=network.target

[Service]
User=andy
WorkingDirectory=/var/www/strava-insight-portal
EnvironmentFile=/var/www/strava-insight-portal/backend/.env
ExecStart=/var/www/strava-insight-portal/venv/bin/uvicorn backend.main:app --host 127.0.0.1 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target
